# iac/ansible/roles/tacacs_server/templates/tacacs-server.service.j2
# SPDX-License-Identifier: Apache-2.0
# TACACS+ Authentication Server - systemd service unit
# Managed by Ansible - do not edit manually

[Unit]
Description=TACACS+ Authentication Server
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/your-org/tacacs-rs

[Service]
Type=notify
User={{ tacacs_user }}
Group={{ tacacs_group }}
ExecStart={{ tacacs_bin_dir }}/tacacs-server \
    --config {{ tacacs_config_dir }}/config.json \
    --log-format {{ tacacs_log_format }} \
{% if tacacs_otel_enabled | bool %}
    --otlp-endpoint {{ tacacs_otel_endpoint }} \
    --otel-service-name {{ tacacs_otel_service_name }} \
    --location {{ tacacs_location }} \
{% endif %}
    --log-level {{ tacacs_log_level }}
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
WatchdogSec=30s
TimeoutStopSec={{ tacacs_drain_timeout_seconds + 30 }}s

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
CapabilityBoundingSet=
AmbientCapabilities=CAP_NET_BIND_SERVICE
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Read-only bind mounts for config
ReadOnlyPaths={{ tacacs_config_dir }}

# Allow writing to log directory
ReadWritePaths={{ tacacs_log_dir }}

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax={{ tacacs_systemd_memory_max }}
CPUQuota={{ tacacs_systemd_cpu_quota }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tacacs-server

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
FinalKillSignal=SIGKILL

[Install]
WantedBy=multi-user.target
