# iac/ansible/roles/pgbouncer/templates/pgbouncer.service.j2
# SPDX-License-Identifier: Apache-2.0
# PgBouncer systemd service unit
# Managed by Ansible - do not edit manually

[Unit]
Description=PgBouncer connection pooler for PostgreSQL
Documentation=https://www.pgbouncer.org/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ pgbouncer_user }}
Group={{ pgbouncer_group }}
ExecStart=/usr/bin/pgbouncer {{ pgbouncer_conf_dir }}/pgbouncer.ini
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# Resource limits
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths={{ pgbouncer_log_dir }} {{ pgbouncer_run_dir }}

# Graceful shutdown
TimeoutStopSec=30s
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
