# iac/packer/files/tacacs-server.service
# SPDX-License-Identifier: AGPL-3.0-only
# TACACS+ Authentication Server - systemd service unit
# Part of golden image, will be customized by Ansible

[Unit]
Description=TACACS+ Authentication Server
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/your-org/tacacs-rs

[Service]
Type=notify
User=tacacs
Group=tacacs
ExecStart=/usr/local/bin/tacacs-server --config /etc/tacacs/config.json
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
WatchdogSec=30s
TimeoutStopSec=60s

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
CapabilityBoundingSet=
AmbientCapabilities=CAP_NET_BIND_SERVICE
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Read-only bind mounts for config
ReadOnlyPaths=/etc/tacacs

# Allow writing to log directory
ReadWritePaths=/var/log/tacacs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tacacs-server

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
FinalKillSignal=SIGKILL

[Install]
WantedBy=multi-user.target
