# iac/ansible/roles/patroni/templates/patroni.service.j2
# SPDX-License-Identifier: Apache-2.0
# Patroni systemd service unit
# Managed by Ansible - do not edit manually

[Unit]
Description=Patroni PostgreSQL HA cluster manager
Documentation=https://patroni.readthedocs.io/
After=network-online.target etcd.service
Wants=network-online.target

[Service]
Type=simple
User={{ patroni_user }}
Group={{ patroni_group }}
ExecStart=/usr/local/bin/patroni /etc/patroni/patroni.yml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10s
TimeoutSec=60

# Environment
Environment=PATRONI_CONFIGURATION=/etc/patroni/patroni.yml

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Security hardening (limited due to PostgreSQL requirements)
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=full
PrivateTmp=yes

# Graceful shutdown
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
