# iac/ansible/roles/frr_anycast/templates/tacacs-anycast-health.service.j2
# SPDX-License-Identifier: Apache-2.0
# TACACS+ Anycast Health Check systemd service
# Managed by Ansible - do not edit manually

[Unit]
Description=TACACS+ Anycast Health Check
Documentation=https://github.com/your-org/tacacs-rs
After=network-online.target frr.service tacacs-server.service
Wants=network-online.target
Requires=frr.service

[Service]
Type=simple
ExecStart=/usr/local/bin/tacacs-health-check.sh
Restart=on-failure
RestartSec=10s

# Environment
Environment=HOME=/root

# Security hardening
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes

# Graceful shutdown - withdraw route before exit
TimeoutStopSec=30s
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
